# Nginx/config-system service
# tell you how to add you nginx service into system service and make the service auto-run when the system reboot
# the code below is undertest
# but the drawback is that you need to make sure the service start first ,otherwise these orders won't work!!!
# add the service into system, so we can use 'systemctl start/stop/reload nginx' to control nginx process! 

# start the service when system start 
# reference is below
$ /etc/rc.d/rc.local 


#!/bin/bash
# chkconfig: - 85 16  #it's very important
# description: Nginx Service Control Script and nginx is a world wide web server!
 case "$1" in
   start)
          /usr/local/nginx/sbin/nginx    
          ;;
   stop)
          /usr/bin/killall -s QUIT nginx
             ;;
   restart)
          $0 stop
          $0 start
          ;;
   reload)
          /usr/bin/killall -s HUP nginx
          ;;
   *)
     echo "Usage:$0 {start|stop|restart|reload}"
     exit 1
 esac
 exit 0
